# kannel-server-cron/build-defaults/etc/supervisor/supervisord.d/kannel.ini
[program:cron]
command = bash -c "while true; chown root:root /etc/cron.d/*; chmod og-wx /etc/cron.d/*; do sleep 0.1; [[ -e /var/run/crond.pid ]] || break; done && exec /usr/sbin/crond -m off -n"
priority = 100
startsecs = 0
autorestart = true
startretries = 9999999999
redirect_stderr = true
stdout_logfile = /var/log/supervisor-logs/cron-daemon.log
stdout_events_enabled = true

[program:bearerbox]
command=/usr/local/kannel/sbin/bearerbox /etc/kannel/kannel.conf
priority=1
process_name=%(program_name)s
redirect_stderr=true
stdout_logfile=/var/log/kannel/supervisor.bearerbox.stdout.log
stderr_logfile=/var/log/kannel/supervisor.bearerbox.stderr.log
user=root
autostart=true
autorestart=true
startretries=2
numprocs=1
stopsignal=QUIT

[program:smsbox]
command=/usr/local/kannel/sbin/smsbox /etc/kannel/kannel.conf
priority=2
process_name=%(program_name)s
redirect_stderr=true
stdout_logfile=/var/log/kannel/supervisor.smsbox.stdout.log
stderr_logfile=/var/log/kannel/supervisor.smsbox.stderr.log
user=root
autostart=true
autorestart=true
startretries=2
numprocs=1
stopsignal=QUIT
