# kannel-server/build-defaults/etc/supervisor/supervisord.d/kannel.ini
[program:bearerbox]
command=/usr/local/kannel/sbin/bearerbox /etc/kannel/kannel.conf
priority=1
process_name=%(program_name)s
redirect_stderr=true
stdout_logfile=/var/log/kannel/supervisor.bearerbox.stdout.log
stderr_logfile=/var/log/kannel/supervisor.bearerbox.stderr.log
user=root
autostart=true
autorestart=true
startretries=2
numprocs=1
stopsignal=QUIT

[program:smsbox]
command=/usr/local/kannel/sbin/smsbox /etc/kannel/kannel.conf
priority=2
process_name=%(program_name)s
redirect_stderr=true
stdout_logfile=/var/log/kannel/supervisor.smsbox.stdout.log
stderr_logfile=/var/log/kannel/supervisor.smsbox.stderr.log
user=root
autostart=true
autorestart=true
startretries=2
numprocs=1
stopsignal=QUIT
